# supervisord.ini
[supervisord]
nodaemon=true
user=root

[program:php]
command=/usr/sbin/php-fpm82 -F
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

